# templates/seafile.service.erb
[Unit]
Description=Daemon to Seafile
After=syslog.target

[Service]
ExecStart=<%= seafile::base_dir %>/seafile-server-latest/seafile.sh restart
Type=forking

[Install]
WantedBy=multi-user.target
